# deploy/bare_metal/systemd/imdf-scoring.service
# ============================================================================
# imdf-scoring.service — aesthetic / quality scoring models (:8005)
# ============================================================================
[Unit]
Description=IMDF scoring-service (aesthetic, CLIP, NSFW scoring)
After=network-online.target postgresql.service redis-server.service
Wants=network-online.target
Requires=postgresql.service redis-server.service

[Service]
Type=simple
User=imdf
Group=imdf
WorkingDirectory=/opt/nanobot-factory
EnvironmentFile=/etc/imdf/imdf.env

ExecStart=/opt/nanobot-factory/venv/bin/uvicorn \
    backend.services.scoring_service.main:app \
    --host 0.0.0.0 \
    --port 8005 \
    --workers 2 \
    --proxy-headers \
    --forwarded-allow-ips='*' \
    --log-level info

Restart=always
RestartSec=5
StartLimitBurst=5
StartLimitIntervalSec=60

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=yes
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictRealtime=true
LockPersonality=true
SystemCallArchitectures=native
ReadWritePaths=/opt/nanobot-factory/data /opt/nanobot-factory/logs

# Scoring may load large ONNX models — give it more RAM
LimitNOFILE=65536
LimitNPROC=4096
MemoryMax=4G
MemoryHigh=3G
CPUQuota=300%

StandardOutput=journal
StandardError=journal
SyslogIdentifier=imdf-scoring

[Install]
WantedBy=multi-user.target