# deploy/bare_metal/systemd/imdf-search.service
# ============================================================================
# imdf-search.service — full-text + vector search (pgvector) (:8011)
# ============================================================================
[Unit]
Description=IMDF search-service (text + vector search)
After=network-online.target postgresql.service redis-server.service
Wants=network-online.target
Requires=postgresql.service redis-server.service

[Service]
Type=simple
User=imdf
Group=imdf
WorkingDirectory=/opt/nanobot-factory
EnvironmentFile=/etc/imdf/imdf.env

ExecStart=/opt/nanobot-factory/venv/bin/uvicorn \
    backend.services.search_service.main:app \
    --host 0.0.0.0 \
    --port 8011 \
    --workers 2 \
    --proxy-headers \
    --forwarded-allow-ips='*' \
    --log-level info

Restart=always
RestartSec=5
StartLimitBurst=5
StartLimitIntervalSec=60

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=yes
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictRealtime=true
LockPersonality=true
SystemCallArchitectures=native
ReadWritePaths=/opt/nanobot-factory/data /opt/nanobot-factory/logs

# Vector index can be large in memory
LimitNOFILE=65536
LimitNPROC=4096
MemoryMax=4G
MemoryHigh=3G
CPUQuota=200%

StandardOutput=journal
StandardError=journal
SyslogIdentifier=imdf-search

[Install]
WantedBy=multi-user.target