# deploy/bare_metal/systemd/imdf-notification.service
# ============================================================================
# imdf-notification.service — email/SMS/webhook dispatch (:8010)
# ============================================================================
[Unit]
Description=IMDF notification-service (email, SMS, webhook)
After=network-online.target postgresql.service redis-server.service
Wants=network-online.target
Requires=postgresql.service redis-server.service

[Service]
Type=simple
User=imdf
Group=imdf
WorkingDirectory=/opt/nanobot-factory
EnvironmentFile=/etc/imdf/imdf.env

ExecStart=/opt/nanobot-factory/venv/bin/uvicorn \
    backend.services.notification_service.main:app \
    --host 0.0.0.0 \
    --port 8010 \
    --workers 2 \
    --proxy-headers \
    --forwarded-allow-ips='*' \
    --log-level info

Restart=always
RestartSec=5
StartLimitBurst=5
StartLimitIntervalSec=60

# Notification dispatch makes outbound HTTPS calls — needs network namespace relax
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=yes
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
LockPersonality=true
SystemCallArchitectures=native
ReadWritePaths=/opt/nanobot-factory/data /opt/nanobot-factory/logs

LimitNOFILE=65536
LimitNPROC=4096
MemoryMax=2G
MemoryHigh=1536M
CPUQuota=150%

StandardOutput=journal
StandardError=journal
SyslogIdentifier=imdf-notification

[Install]
WantedBy=multi-user.target